# _docker/app/crontab/root
# do daily/weekly/monthly maintenance
# min	hour	day	month	weekday	command
*/1   *   *   *   *   echo "This is test" >> /var/log/cron/cron.log 2>&1